# debian/storiqone-daemon.postinst
#!/bin/sh

set -e

if [ "$1" = configure ]; then
	chkconfig stoned on
fi
